# t/fixtures/runnable/crontab
*/10 * * * * echo "cron" > /tmp/cron1
*/10 * * * * echo "cron2" > /tmp/cron1
0 09-18 * * 1-5 echo "cron3" > /tmp/cron1
