# t/fixtures/runnable/crontab
*/1 * * * * /usr/app/t/fixtures/script.sh min
*/2 * * * * /usr/app/t/fixtures/script.sh 2min
*/10 * * * * /usr/app/t/fixtures/script.sh 10min
*/15 * * * * /usr/app/t/fixtures/script.sh 15min
